# hoot-jetty/src/osx/bin/rc.local
#!/bin/sh -e
# Copy this line to your rc.local or contribute proper init scripts :).
/bin/su - hoot -c /home/hoot/jetty/start.sh
